# deploy/reel2bits.target
[Unit]
Description=reel2bits
Wants=reel2bits-web.service reel2bits-worker.service
